# crontab
MAILTO=user@example.org

*/5 * * * * root path/to/command
